# scripts/dajtu-backup.timer
[Unit]
Description=Hourly dajtu backup

[Timer]
OnCalendar=*:11:05
Persistent=true

[Install]
WantedBy=timers.target
